# deploy/systemd/podcast-schedule-tick.service
[Unit]
Description=Podcast Hub schedule tick
After=network-online.target podcast-api.service

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -fsS -X POST http://127.0.0.1:14000/internal/schedules/tick

